# zta-windows/playbooks/templates/nmap.service.j2
[Unit]
Description=nmap

[Service]
ExecStart=-/usr/bin/nmap  -p 8082 {{ hostvars['windows-ad']['private_ip2'] }}